# tbtacl/tbtacl.rules.in
# Thunderbolt udev rules for ACL (device auto approval)
SUBSYSTEM=="thunderbolt" ENV{DEVTYPE}=="thunderbolt_device" ACTION=="add"    ATTR{authorized}=="0" RUN+="@UDEV_BIN_DIR@/tbtacl add    $devpath"
SUBSYSTEM=="thunderbolt" ENV{DEVTYPE}=="thunderbolt_device" ACTION=="change" ATTR{authorized}!="0" RUN+="@UDEV_BIN_DIR@/tbtacl change $devpath"
